# etc/systemd/system/secondarynamenode.service
[Unit]
Description=HDFS Secondary-namenode Daemon
Documentation=http://hadoop.apache.org
Requires=network.target
After=network.target
 
[Service]   
Environment="JAVA_HOME=/usr/lib/jvm/java-1.8.0-openjdk"
Environment="HADOOP_HOME=/usr/local/hadoop"
Environment="HADOOP_CONF_DIR=/etc/hadoop/conf"
Type=forking
WorkingDirectory=/usr/local/hadoop
User=hdfs
Group=hadoop
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p /var/run/hadoop/hdfs/
ExecStartPre=/bin/chown -R hdfs:hadoop /var/run/hadoop/hdfs/
ExecStart=/usr/local/hadoop/bin/hdfs --daemon start secondarynamenode
ExecStop=/usr/local/hadoop/bin/hdfs --daemon stop secondarynamenode
ExecReload=/usr/local/hadoop/bin/hdfs --daemon restart secondarynamenode
TimeoutSec=30
Restart=on-failure
 
[Install]
WantedBy=default.target

